# files/linux_filesystems/media-films.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=remote-fs.target

[Mount]
What=kalibur-mce:movies-vol
Where=/media/films
Type=glusterfs
Options=rw,relatime,user_id=0,group_id=0,allow_other,max_read=131072,_netdev,log-level=WARNING,log-file=/var/log/movies-vol-gluster.log,backup-volfile-servers=10.42.0.2:kal-bunker:10.42.0.30:/movies-vol
