# files/linux_filesystems/media-incoming.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=remote-fs.target

[Mount]
What=kalibur-mce:/media/incoming
Where=/media/incoming
Type=nfs
Options=_netdev,soft,defaults
